# xray.initd
#!/sbin/openrc-run

name="Xray"
description="A unified platform for anti-censorship."
supervisor="supervise-daemon"
command="/usr/local/bin/xray"
command_args="run -c /usr/local/etc/xray/config.json"

depend() {
    need net
}